# kak/.config/profile.d/10-kak.sh
export EDITOR=kak
export SUDO_EDITOR=kak
