# npm/.config/profile.d/50-npm.sh
export PATH="$HOME/.local/share/npm/bin:$PATH"
