# roarke.host/.config/profile.d/10-roarke.sh
export EDAC_PULSE_NAME=alsa_card.usb-FiiO_FiiO_USB_DAC-E18-01 
export EDAC_ALSA_NAME=hw:DACE18
