# sway/.config/profile.d/99-sway.sh
if [ "$(tty)" = "/dev/tty1" ]; then
    sway
    exit 0
fi
