# go/.config/profile.d/50-go.sh
GOPATH="$HOME/.local/go"
export GOPATH
PATH="$GOPATH/bin:$PATH"
export PATH
