# mail/.config/profile.d/50-notmuch.sh
NOTMUCH_CONFIG="$HOME/.config/notmuch/config"
export NOTMUCH_CONFIG
