# ssh/.config/profile.d/10-ssh.sh
export SSH_ASKPASS=/usr/lib/ssh/gnome-ssh-askpass2
eval $(keychain --eval --noask --quiet) >/dev/null 2>/dev/null
