# ssh/.config/profile.d/10-ssh.sh
export SSH_ASKPASS=/usr/lib/ssh/ssh-askpass
eval $(keychain --eval --noask --quiet)
